# buildroot/rootfs_overlay/etc/profile.d/log_login.sh
if [ "$(id -u)" != "0" ]; then 
	logger -p auth.info "$(whoami) login on 'console'"
fi
